# docker/backup/crontab
0 2 * * * root /usr/local/bin/entrypoint.sh
